# br2-external/package/trezord-go/trezord.service
[Unit]
Description=Trezor Bridge Daemon
Documentation=https://github.com/trezor/trezord-go
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/trezord -e 21324
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
User=root
Group=root
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateNetwork=false
PrivateDevices=false

# Allow access to USB devices
SupplementaryGroups=plugdev

# Logging
SyslogIdentifier=trezord

[Install]
WantedBy=multi-user.target
